# deploy/enrichment-mcp.service
[Unit]
Description=Enrichment MCP Server (FastMCP, Streamable HTTP)
Documentation=https://github.com/jlowin/fastmcp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Non-root service user -- create with:
#   useradd --system --no-create-home --shell /usr/sbin/nologin mcpsvr
User=mcpsvr
WorkingDirectory=/opt/enrichment-mcp
EnvironmentFile=/opt/enrichment-mcp/.env

# FastMCP server entry point; reads MCP_HOST and MCP_PORT from .env
ExecStart=/opt/enrichment-mcp/.venv/bin/python -m src.mcp_server

Restart=on-failure
RestartSec=10

# Limit the blast radius if something goes wrong
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/enrichment-mcp

[Install]
WantedBy=multi-user.target
